# system/files/xinetd.d/Makefile
#!/usr/bin/make -f

include ../../configure.mk

install:
	./install-listener.sh install $(INST_SBINDIR)
uninstall:
	./install-listener.sh uninstall $(INST_SBINDIR)
